# crontab
0 1 * * * python3 DEV_EOD_Stock_Email.ipynb >/proc/1/fd/1 2>/proc/1/fd/2

# An empty line is required at the end of this file for a valid cron file.
